# robot-os/arduino/ros_arduino-indigo-devel/ros_arduino_base/devel/etc/catkin/profile.d/05.catkin-test-results.sh
# generated from catkin/cmake/env-hooks/05.catkin-test-results.sh.develspace.in

export CATKIN_TEST_RESULTS_DIR="/home/ubuntu/p/compv/robot-os/arduino/ros_arduino-indigo-devel/ros_arduino_base/test_results"
export ROS_TEST_RESULTS_DIR="$CATKIN_TEST_RESULTS_DIR"
